# infra/systemd/loop-farm-control.service
[Unit]
Description=Loop Farm Control API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/loop-farm
Environment=LOOP_FARM_HOST=0.0.0.0
Environment=LOOP_FARM_PORT=8787
Environment=LOOP_FARM_DB=/var/lib/loop-farm/control.sqlite3
Environment=LOOP_FARM_ADMIN_TOKEN=change-me
ExecStart=/opt/loop-farm/.venv/bin/loop-farm-control
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

